# debian/prerm
systemctl stop crowdsec-aws-waf-bouncer || echo "cannot stop service"
systemctl disable crowdsec-aws-waf-bouncer || echo "cannot disable service"